# HW6-Systemd/servicekey/servicekey.service
[Unit]
Description=Script for finding keyword

[Service]
Type=oneshot
EnvironmentFile=/etc/sysconfig/servicekey
ExecStart=/usr/bin/servicekey.sh